# debian/paquet-openldap-debian/openldap-ltb-2.4.50/debian/openldap-ltb.prerm
#!/bin/bash

set -e

BASEDIR=$(dirname $0)
source $BASEDIR/openldap-ltb.vars


# stops OpenLDAP and removes startup links
if [ -x "/etc/init.d/slapd" ]; then
        if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
                invoke-rc.d slapd stop || exit $?
        else
                /etc/init.d/slapd stop || exit $?
        fi
fi

